# dot_config/systemd/user/ssh-agent.service
# Shell-independent ssh-agent: no shell rc starts an agent (the bash →
# zsh switch on Omarchy silently dropped it), so the user session owns
# it instead. SSH_AUTH_SOCK reaches shells via the guarded fallback in
# .zshrc/.bashrc — deliberately NOT via environment.d, which would force
# the socket on the whole session even when this unit isn't running.
# Enabled automatically by .chezmoiscripts/run_onchange_after_10-enable-ssh-agent.sh.tmpl
[Unit]
Description=SSH key agent

[Service]
Type=simple
Environment=SSH_AUTH_SOCK=%t/ssh-agent.socket
ExecStart=/usr/bin/ssh-agent -D -a $SSH_AUTH_SOCK

[Install]
WantedBy=default.target
